# sys-cluster/hadoop-bin/files/hadoop-namenode.service
[Unit]
Description=Hadoop Namenode
After=network.target

[Service]
EnvironmentFile=/etc/conf.d/hadoop
ExecStart=/usr/bin/hdfs namenode
User=hdfs
Group=hadoop

[Install]
WantedBy=multi-user.target
